# settings/site.conf
[program:site]
command=/dmojsite/bin/uwsgi --ini uwsgi.ini
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log